# templates/default/systemd.service.erb
[Unit]
Description=<%= @description %>

[Service]
Type=simple
Restart=on-failure
ExecStart=<%= @daemon %>
<% if @env %>
Environment="<%= @env %>"
<% end -%>
<% if @user %>
User=<%= @user %>
<% end -%>
<% if @chdir %>
WorkingDirectory=<%= @chdir %>
<% end -%>

[Install]
WantedBy=multi-user.target
